# BALINTONI_BENIAMIN/PoriectDATC/Raspberry/datc.service
[Unit]
Description=My DATC 

[Service]
ExecStart=/usr/bin/Temp

[Install]
WantedBy=multi-user.target
